# root/etc/s6-overlay/s6-rc.d/svc-mod-uds-proxy/run
#!/usr/bin/with-contenv bash

# echo "starting nginx" && nginx -t && nginx -c /default/nginx.conf -g 'daemon off;'
