# root/etc/s6-overlay/s6-rc.d/svc-mod-uds-proxy/run
#!/usr/bin/with-contenv bash

echo "[UDS_PROXY] starting nginx" && nginx -g 'daemon off;'
