# root/etc/s6-overlay/s6-rc.d/init-mod-uds-proxy-add-package/run
#!/usr/bin/with-contenv bash

echo "[UDS_PROXY] adding nginx to the package install list"

echo "nginx" >> /mod-repo-packages-to-install.list