# root/etc/s6-overlay/s6-rc.d/init-mod-uds-proxy/run
#!/usr/bin/with-contenv bash

## Ubuntu
if [ -f /usr/bin/apt ]; then
    echo "ubuntu"

fi
# Alpine
if [ -f /sbin/apk ]; then
    echo "alpine"

echo "nginx" >> /mod-repo-packages-to-install.list
